# deploy/aria-heartbeat.timer.in
[Unit]
Description=Run ARIA heartbeat every 5 minutes (external watchdog)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
